# deploy/conf.d/mapainteractivo_dev.ini
[program:mapainteractivo-web]
command = /apps/mapainteractivo/deploy/scripts/web_dev.sh
user = nginx
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile_maxbytes = 1MB
stdout_logfile = /apps/mapainteractivo/deploy/logs/web.log

[group:mapainteractivo]
programs = mapainteractivo-web
